# supabase/functions/_cron/daily.cron
# Run daily at 3:00 AM UTC
0 3 * * * daily-cron
